# scripts/copy_grid_proxy/copy-grid-proxy.service
[Unit]
Description=Service to copy the x509 proxy to the XRootD buffer

[Service]
Type=oneshot
ExecStart=/path/to/copy_grid_proxy.sh
