# scripts/grid-security/gridsecurity.timer
[Unit]
Description=Update grid-security every 24h

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

