# scripts/grid-security/gridsecurity-restart.service
[Unit]
Description=Force-restart gridsecurity.service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart gridsecurity.service